# piip.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/piip.sh

[Install]
WantedBy=default.target